# deploy/com.mail-listener.gmail-imap-idle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.mail-listener.gmail-imap-idle</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/env</string>
      <string>python3</string>
      <string>-m</string>
      <string>mail_listener</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>REPLACE_WITH_PROJECT_ROOT</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>MAIL_LISTENER_BASE_DIR</key>
      <string>/Users/REPLACE_ME/.cache/mail-listener</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/REPLACE_ME/.cache/mail-listener/logs/launchd.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_ME/.cache/mail-listener/logs/launchd.stderr.log</string>
  </dict>
</plist>
